# logiclemonai.service
[Unit]
Description=Logiclemonai Secure Connect Server
After=network.target

[Service]
Type=simple
User=toofargonetrading
WorkingDirectory=/home/toofargonetrading/logiclemonai
ExecStart=/usr/bin/node /home/toofargonetrading/logiclemonai/secure-connect-landing/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target